# lib/.config/profile.d/keyboard
#!/usr/bin/env sh

setxkbmap -layout us -variant altgr-intl || true
