# lib/.config/profile.d/program_settings
#!/usr/bin/env sh

export DICS="/usr/share/stardict/dic/"
export DOTNET_CLI_TELEMETRY_OPTOUT=1
export SUDO_ASKPASS="$HOME/.local/bin/dmenupass"
export FZF_DEFAULT_OPTS="--layout=reverse --height 40%"
export LESS=-R
export LESS_TERMCAP_mb="$(printf '%b' '[1;31m')"
export LESS_TERMCAP_md="$(printf '%b' '[1;36m')"
export LESS_TERMCAP_me="$(printf '%b' '[0m')"
export LESS_TERMCAP_so="$(printf '%b' '[01;44;33m')"
export LESS_TERMCAP_se="$(printf '%b' '[0m')"
export LESS_TERMCAP_us="$(printf '%b' '[1;32m')"
export LESS_TERMCAP_ue="$(printf '%b' '[0m')"
export QT_QPA_PLATFORMTHEME=gtk3
export DOTFILES="$HOME/doc/src/github.com/kruhlmann/dotfiles"
export GCC_COLORS='error=01;31:warning=01;35:note=01;36:caret=01;32:locus=01:quote=01'
export GITLAB_TOKEN="$(cat "$HOME"/.repo-cache)"
export GITLAB_ENDPOINT="https://mo-gitlab.siemens.dk/api/v4"
export GITLAB_REPO_CACHE="$HOME/.repo-cache"
export LF_ICONS="di=ï“:fi=ï:ln=ï’:or=ï’:ex=ï§:*.c=î˜:*.cc=î˜:*.clj=î¨:*.coffee=î‘:*.cpp=î˜:*.css=î˜”:*.d=î¯:*.dart=î˜:*.erl=î±:*.exs=î˜­:*.fs=î§:*.go=î˜§:*.h=î˜:*.hh=î˜:*.hpp=î˜:*.hs=î·:*.html=î˜:*.java=îœ¸:*.jl=î˜¤:*.js=î:*.json=î˜‹:*.lua=î˜ :*.md=î˜‰:*.php=î˜ˆ:*.pl=î©:*.pro=î¡:*.py=îœ¼:*.rb=îœ¹:*.rs=î¨:*.scala=îœ·:*.ts=î˜¨:*.vim=î˜«:*.cmd=ï:*.ps1=ï:*.sh=ï:*.bash=ï:*.zsh=ï:*.fish=ï:*.tar=ï:*.tgz=ï:*.arc=ï:*.arj=ï:*.taz=ï:*.lha=ï:*.lz4=ï:*.lzh=ï:*.lzma=ï:*.tlz=ï:*.txz=ï:*.tzo=ï:*.t7z=ï:*.zip=ï:*.z=ï:*.dz=ï:*.gz=ï:*.lrz=ï:*.lz=ï:*.lzo=ï:*.xz=ï:*.zst=ï:*.tzst=ï:*.bz2=ï:*.bz=ï:*.tbz=ï:*.tbz2=ï:*.tz=ï:*.deb=ï:*.rpm=ï:*.jar=ï:*.war=ï:*.ear=ï:*.sar=ï:*.rar=ï:*.alz=ï:*.ace=ï:*.zoo=ï:*.cpio=ï:*.7z=ï:*.rz=ï:*.cab=ï:*.wim=ï:*.swm=ï:*.dwm=ï:*.esd=ï:*.jpg=ï:*.jpeg=ï:*.mjpg=ï:*.mjpeg=ï:*.gif=ï:*.bmp=ï:*.pbm=ï:*.pgm=ï:*.ppm=ï:*.tga=ï:*.xbm=ï:*.xpm=ï:*.tif=ï:*.tiff=ï:*.png=ï:*.svg=ï:*.svgz=ï:*.mng=ï:*.pcx=ï:*.mov=ï:*.mpg=ï:*.mpeg=ï:*.m2v=ï:*.mkv=ï:*.webm=ï:*.ogm=ï:*.mp4=ï:*.m4v=ï:*.mp4v=ï:*.vob=ï:*.qt=ï:*.nuv=ï:*.wmv=ï:*.asf=ï:*.rm=ï:*.rmvb=ï:*.flc=ï:*.avi=ï:*.fli=ï:*.flv=ï:*.gl=ï:*.dl=ï:*.xcf=ï:*.xwd=ï:*.yuv=ï:*.cgm=ï:*.emf=ï:*.ogv=ï:*.ogx=ï:*.aac=ï:*.au=ï:*.flac=ï:*.m4a=ï:*.mid=ï:*.midi=ï:*.mka=ï:*.mp3=ï:*.mpc=ï:*.ogg=ï:*.ra=ï:*.wav=ï:*.oga=ï:*.opus=ï:*.spx=ï:*.xspf=ï:*.pdf=ï‘"
export SONAR_SCANNER_HOME="/opt/sonar-scanner"
export FZF_DEFAULT_COMMAND='fd --type f --hidden --follow --exclude .git'
export FZF_CTRL_T_COMMAND="$FZF_DEFAULT_COMMAND"
export EDITOR=/usr/bin/nvim
export SSH_AUTH_SOCK="$(gpgconf --list-dirs agent-ssh-socket)"
