# lib/.config/systemd/user/synccal.service
[Unit]
Description=Synchronize calendar with remote git repo

[Service]
ExecStart=/home/ges/.scripts/synccal
