# lib/.config/systemd/user/sync_music.timer
[Unit]
Description=Synchronize music with YouTube playlists

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min 

[Install]
WantedBy=timers.target
