# lib/.config/systemd/user/sync_music.service
[Unit]
Description=Synchronize music with YouTube playlists 

[Service]
ExecStart=%h/.scripts/sync_music
