# lib/.config/systemd/user/syncmail.service
[Unit]
Description=Synchronize Siemens email

[Service]
ExecStart=/home/ges/.scripts/syncmail
